# horizon/files/supervisor-horizon.conf
[program:horizon]
command=python manage.py runserver
directory=/opt/openstack/horizon/
