# keystone/files/supervisor-keystone.conf
[program:keystone]
user=keystone
command=/usr/local/bin/keystone-all
