# glance/files/supervisor-glance.conf
[program:glance-api]
user=glance
command=/usr/local/bin/glance-api

[program:glance-registry]
user=glance
command=/usr/local/bin/glance-registry
